# infra/crs-copy-logs.timer
[Unit]
Description=Run copy every 15 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min
Unit=crs-copy-logs.service

[Install]
WantedBy=timers.target