# infra/crs-copy-logs.service
[Unit]
Description=Copy off logs and data

[Service]
Type=oneshot
User=crs
Group=crs
WorkingDirectory=/crs
EnvironmentFile=/etc/environment
Environment=AZCOPY_CONCURRENCY_VALUE=1000
Environment=GOMAXPROCS=2

ExecStartPre=-+zfs destroy data@azcopy
ExecStartPre=+zfs snapshot data@azcopy
ExecStart=azcopy sync --log-level=NONE /data/.zfs/snapshot/azcopy/ "${CRS_BLOB_ENDPOINT}/backup/" --exclude-path="docker;influxdb;tmp"
ExecStartPost=-+zfs destroy data@azcopy
