# infra/crs-task-server.service
[Unit]
Description=CRS Task Server
After=network.target

[Service]
Type=simple
User=crs
Group=crs
WorkingDirectory=/crs
EnvironmentFile=/etc/environment
Environment=CRS_SERVICE_NAME=roboduck.task_server
ExecStart=/crs/.venv/bin/python3 -m uvicorn crs.task_server.app:app --port 1324
Restart=always

[Install]
WantedBy=multi-user.target
